# assets/linux/service/php85-fpm.service
cat > /etc/systemd/system/php85-fpm.service << 'EOF'
[Unit]
Description=The PHP 8.5 FastCGI Process Manager
After=network.target

[Service]
Type=notify
UMask=0027
User=php
Group=php

RuntimeDirectory=php85-fpm
RuntimeDirectoryMode=0750

ExecStart=/server/php/85/sbin/php-fpm --nodaemonize --fpm-config /server/etc/php/85/php-fpm.conf
ExecReload=/bin/kill -USR2 $MAINPID
ExecStop=/bin/kill -SIGQUIT $MAINPID

KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=90

Restart=on-failure
RestartSec=10s

LimitNOFILE=65536

ProcessMax=1024

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_NETLINK AF_UNIX

ReadWritePaths=/server/www /server/logs/php

[Install]
WantedBy=multi-user.target
EOF
