# assets/linux/service/redis.service
cat > /etc/systemd/system/redis.service << 'EOF'
[Unit]
Description=Redis In-Memory Data Store
After=network.target

[Service]
Type=notify
UMask=0027
User=redis
Group=redis

RuntimeDirectory=redis
RuntimeDirectoryMode=0750

ExecStart=/server/redis/bin/redis-server /server/etc/redis/config/redis.conf
ExecReload=/bin/kill -s HUP $MAINPID

KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=90

Restart=on-failure
RestartSec=10s

LimitNOFILE=10000
LimitNPROC=10000
LimitMEMLOCK=infinity
LimitAS=infinity

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_NETLINK AF_UNIX

ReadWritePaths=/server/logs/redis

[Install]
WantedBy=multi-user.target
EOF
